# helloworld.conf
description "Hello world Daemon"

# Start when the system is ready to do networking.
start on started elastic-network-interfaces

# Stop when the system is on its way down.
stop on shutdown

# Automatically restart the daemon if it crashes
respawn

script
    # Execute the management script
    exec /home/ec2-user/helloworld-manager.sh
end script
